# lab5/crontab
* * * * 1-5 bash /vagrant/task_B.sh > /var/log/logged-in.log 
