# lab5/daemon.service
[Unit]
Description=DaemonDevOps

[Service]
ExecStart=/usr/bin/daemon/daemon.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target